# host-files/papermc-taijotain-restic-backup.timer
[Unit]
Description=Restic backup of the PaperMC data (taijotain)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
